# Carleton_ROV_File_System/Scripts/UART1_SERVICE.service

[Unit]
Description=THIS IS THE UART1 INPUT/OUPUT SERVICE 

After=syslog.target network.target

[Service]
Type=simple
ExecStart=/usr/bin/UART1_XBEE.sh

[Install]
WantedBy=multi-user.target